# infra/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:backend]
command=node /app/server/server.js
directory=/app/server
autostart=true
autorestart=true
environment=NODE_ENV=production,PORT=3001

[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile
autostart=true
autorestart=true
